# v0.1/storm/graph/storm/files/nimbus.cfg
[supervisord]
logfile=/home/log/super/supervisord.log
logfile_maxbytes=20MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid

[program:zookeeper]
command=/home/zookeeper-3.4.5/bin/zkServer.sh start-foreground
autostart=true
autorestart=true
log_stdout=true
log_stderr=true
logfile=/home/log/super/zookeeper.log
logfile_maxbytes=20MB
logfile_backups=10

[program:storm-nimbus]
command=storm nimbus
autostart=true
autorestart=true
log_stdout=true
log_stderr=true
logfile=/home/log/super/nimbus.log
logfile_maxbytes=20MB
logfile_backups=10

[program:storm-ui]
command=storm ui
autostart=true
autorestart=true
log_stdout=true
log_stderr=true
logfile=/home/log/super/ui.log
logfile_maxbytes=20MB
logfile_backups=10

