# storm-bash/src/main/model/graph/storm/files/worker.cfg
[supervisord]
logfile=/home/log/super/supervisord.log
logfile_maxbytes=20MB        
logfile_backups=10
loglevel=info                
pidfile=/tmp/supervisord.pid 

[program:storm-supervisor]
command=storm supervisor
autostart=true
autorestart=true
log_stdout=true
log_stderr=true
logfile=/home/log/super/supervisor.log
logfile_maxbytes=20MB
logfile_backups=10

