# storm-bash/src/main/model/graph/storm_nimbus/files/nimbus.cfg
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/var/log/supervisord/supervisord.log
logfile_maxbytes=20MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl] 
serverurl=unix:///tmp/supervisor.sock

[program:zookeeper]
command=/usr/local/zookeeper/bin/zkServer.sh start-foreground
autostart=true
autorestart=true
log_stdout=true
log_stderr=true
logfile=/var/log/zookeeper/zookeeper.log
logfile_maxbytes=20MB
logfile_backups=10

[program:storm-nimbus]
command=storm nimbus
autostart=true
autorestart=true
log_stdout=true
log_stderr=true
logfile=/var/log/storm/nimbus.log
logfile_maxbytes=20MB
logfile_backups=10

[program:storm-ui]
command=storm ui
autostart=true
autorestart=true
log_stdout=true
log_stderr=true
logfile=/var/log/storm/ui.log
logfile_maxbytes=20MB
logfile_backups=10

